# google_dyndns.cron
*/15 * * * * /opt/google_dyndns/.venv/bin/python /opt/google_dyndns/google_dyndns.py
